# etc/services.d/geoip/run
#!/usr/bin/with-contenv sh
exec /bin/bash /xshok-geoip.sh
